# cron/com.agentrealestate.freshness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- staged, NOT loaded. 활성화는 사용자 명시 opt-in (agent_money 와 동일 정책).
     cp cron/com.agentrealestate.freshness.plist ~/Library/LaunchAgents/
     launchctl load ~/Library/LaunchAgents/com.agentrealestate.freshness.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.agentrealestate.freshness</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Volumes/EXT_SSD/bot/agent_realestate/cron/run_agent_realestate.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict><key>Day</key><integer>1</integer><key>Hour</key><integer>9</integer></dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>/Volumes/EXT_SSD/bot/agent_realestate/data/cron.log</string>
  <key>StandardErrorPath</key><string>/Volumes/EXT_SSD/bot/agent_realestate/data/cron.err</string>
</dict>
</plist>
